# systemd/grafana-password-refresh.service
[Unit]
Description=Refresh Grafana admin password from SSM Parameter Store
After=network-online.target docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/refresh-grafana-password.sh
