# systemd/prometheus-creds-refresh.service
[Unit]
Description=Refresh EC2 credentials for Prometheus ec2_sd
After=network-online.target docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/refresh-ec2-credentials.sh
